# Projects/ServerManager/server-monitor-package/systemd/server-monitor.service
[Unit]
Description=Server Monitor Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/server-monitor
ExecStart=/usr/bin/python3 /opt/server-monitor/server_monitor.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
